# cloudbuild/supervisord.conf
[supervisord]
nodaemon=true
[program:cron]
command=crond
